# docker_flask/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
username = dummy
password = dummy

[inet_http_server]
port=127.0.0.1:9001
username = dummy
password = dummy

[supervisord]
nodaemon=true
user=root
loglevel=warn

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock
username = dummy
password = dummy

[eventlistener:flush_signalbus_trigger]
command=/usr/src/app/docker_flask/trigger_supervisor_process.py signalbus_flush 60
directory=/usr/src/app
events=TICK_5

[program:signalbus_flush]
command=/usr/src/app/docker_flask/flush_signalbus.py
directory=/usr/src/app
autostart=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr=true
startsecs=0
exitcodes=0
autorestart=false
